# etc/xdg/autostart/welcome-studio.desktop
[Desktop Entry]
Exec=welcome-studio
Icon=welcome-studio
Name=Welcome  Studio TigerOS
StartupNotify=true
Terminal=false
Type=Application